# phc2sys/phc2sys.service
[Unit]
Description=Synchronize system clock to PTP hardware clock
After=ts2phc-gps.service
Requires=ts2phc-gps.service

[Service]
ExecStart=/usr/sbin/phc2sys -r -s /dev/ptp1 -c CLOCK_REALTIME -w
CPUSchedulingPolicy=rr
CPUSchedulingPriority=98
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
